# crontab
0 */3 * * * /home/jnakano/tmp/ts2mp4/scripts/rm_old_finished.sh
